# deploy/systemd/edx-update.timer
# e-disclosure extractor — systemd timer.
#
# OnCalendar=*-*-* 04:00:00 mirrors config/app.yaml → schedule.cron_time.
# Persistent=true catches up if the host was offline at the scheduled time
# (e.g. laptop wake-from-sleep) without queueing the run twice.

[Unit]
Description=Daily 04:00 trigger for edx-update.service
Requires=edx-update.service

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=10min
Unit=edx-update.service

[Install]
WantedBy=timers.target
