# deploy/systemd/edx-evolve.service
# Patch 44: self-evolve tick runner. Triggered every 5 minutes by
# edx-evolve.timer. Serialised against itself AND against the daily
# edx-update.service via flock on /tmp/edx-evolve.lock.
#
# Failures of the underlying ETL (per-company) do NOT make the
# service unit fail — only crashes (non-zero exit + traceback) do.
# This is intentional: systemd shouldn't go red because of a normal
# failed tick.
#
# Install with:
#   sudo cp deploy/systemd/edx-evolve.service /etc/systemd/system/
#   sudo cp deploy/systemd/edx-evolve.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now edx-evolve.timer

[Unit]
Description=edx self-evolve tick (one batch of 3 companies)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/edx
EnvironmentFile=-/opt/edx/.env
EnvironmentFile=-/opt/edx/.env.evolve
# Proxy env for hosts where direct egress to api.anthropic.com is blocked
# (Anthropic returns 403 "Request not allowed" for direct requests from
# such regions). systemd does NOT source ~/.bashrc, so proxy vars set
# there never reach the agent subprocess. Operator owns this file:
# /opt/edx/.env.proxy (chmod 600), with HTTPS_PROXY/HTTP_PROXY/NO_PROXY.
EnvironmentFile=-/opt/edx/.env.proxy
ExecStart=/usr/bin/flock -n /tmp/edx-evolve.lock /opt/edx/.venv/bin/edx evolve tick
User=edx
Group=edx
TimeoutStartSec=13h
# Allow the npm-installed claude binary on PATH.
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
